# sport-event-bot.service
[Unit]
Description=Telegram Sport Event Bot
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/usr/local/tgbot
EnvironmentFile=/usr/local/tgbot/.env
ExecStart=/usr/local/tgbot/run_sport_event_bot.sh
Restart=always
RestartSec=10
StandardOutput=append:/usr/local/tgbot/sport_event_bot/logs/systemd.log
StandardError=append:/usr/local/tgbot/sport_event_bot/logs/systemd.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
